# duplicati-canary-bin/duplicati.service
[Unit]
Description=Duplicati
After=network.target

[Service]
ExecStart=/opt/duplicati/duplicati-server $DAEMON_OPTS
Restart=on-abort
EnvironmentFile=-/etc/default/duplicati
User=duplicati
Group=duplicati

[Install]
WantedBy=default.target
